# deleteblock.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
HOME=/develop/btc2
0 0  * * * root python /develop/btc2/delete.py
