# syncblock.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
HOME=/develop/btc2
*/5 *  * * * root python /develop/btc2/block.py
